# scripts/cron/com.mir.archive_detector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mir.archive_detector</string>
  <key>ProgramArguments</key>
  <array>
    <string><your-harness-path>/.venv/bin/python</string>
    <string>-m</string>
    <string>tools.fleet_observe.archive.detector</string>
    <string>--digest</string>
    <string>--format</string>
    <string>markdown</string>
  </array>
  <key>WorkingDirectory</key>
  <string><your-harness-path></string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>0</integer>
    <key>Hour</key>
    <integer>6</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/mir-archive-detector-weekly.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/mir-archive-detector-weekly.err</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
