# scripts/cron/com.mir.memory_gc_daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mir.memory_gc_daily</string>
  <key>ProgramArguments</key>
  <array>
    <string><your-harness-path>/.venv/bin/python</string>
    <string><your-harness-path>/scripts/memory_gc_runner.py</string>
    <string>--dry-run</string>
  </array>
  <key>WorkingDirectory</key>
  <string><your-harness-path></string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>4</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/mir-memory-gc.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/mir-memory-gc.err</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
