# scripts/cron/com.mir.template_health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mir.template_health</string>
  <key>ProgramArguments</key>
  <array>
    <string><your-harness-path>/.venv/bin/python</string>
    <string>-m</string>
    <string>tools.fleet_observe.template_health</string>
    <string>--template-repo</string>
    <string><your-harness-path><example-family></string>
    <string>--families-from</string>
    <string><your-harness-path>/config/fleet-harness-state.json</string>
    <string>--output</string>
    <string>json</string>
    <string>--report</string>
    <string>/tmp/mir-template-health-daily.json</string>
  </array>
  <key>WorkingDirectory</key>
  <string><your-harness-path></string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>13</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/mir-template-health.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/mir-template-health.err</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
